# borg/backup.timer
[Unit]
Description=Borg Backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
