# ttrss/update-ttrss.timer
[Unit]
Description=Update Tiny Tiny RSS

[Timer]
OnCalendar=12:00
Persistent=true

[Install]
WantedBy=timers.target
