# iptables/firewall-reload.path
[Unit]
Description=Reload firewall rules when changed

[Path]
PathChanged=/usr/local/etc/firewall/

[Install]
WantedBy=multi-user.target
